# templates/dsid.service.epp
<%- | String $env_file,
      String $current_path,
| -%>
[Unit]
Description=Distributed Service Interface Daemon (DSID)
Requires=network-online.target
After=network-online.target

[Service]
EnvironmentFile=<%= $env_file %>
User=dsid
Group=daq
Type=forking
ExecStart=<%= $current_path %>/x86/bin/dsid_startup ${backingdir} ${interface}
#Restart=always
#RestartSec=10
CPUSchedulingPolicy=fifo

[Install]
WantedBy=multi-user.target
